# service/opendkim/run
#!/bin/sh
exec 2>&1
exec /usr/sbin/opendkim -f -x /mail/etc/opendkim/opendkim.conf