# etc/crontab
@hourly /sbin/sv hup /service/postfix
@hourly /sbin/sv hup /service/dovecot
